# setup/systemd/ss-audioplayback.service


[Unit]
Description=Audio playback service
Wants=ss-inf.service
After=ss-inf.service

[Service]
User=pi
Type=simple
ExecStart=ss-audioplayback.sh

[Install]
WantedBy=multi-user.target

