# setup/systemd/ss-notif-baby_cry.service
[Unit]
Description=Notification (Baby Crying) Service
Wants=ss-yamnet-inf.service
After=ss-yamnet-inf.service

[Service]
User=pi
Type=simple
ExecStart=ss-notification.sh 20 .4

[Install]
WantedBy=multi-user.target

