# setup/systemd/ss-notif-conversation.service
[Unit]
Description=Notification (Conversation) Service
Wants=ss-yamnet-inf.service
After=ss-yamnet-inf.service

[Service]
User=pi
Type=simple
ExecStart=ss-notification.sh 2 .4

[Install]
WantedBy=multi-user.target

